# etc/rc.d/altqd
#!/bin/sh
#
# $NetBSD: altqd,v 1.3 2003/03/10 18:17:06 lukem Exp $
#

# PROVIDE: altqd
# REQUIRE: mountcritremote
# BEFORE:  SERVERS
# KEYWORD: shutdown

. /etc/rc.subr

name="altqd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/altq.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
