# etc/rc.d/amd
#!/bin/sh
#
# $NetBSD: amd,v 1.2 2000/04/15 21:14:50 tsarna Exp $
#

# PROVIDE: amd
# REQUIRE: portmap mountall ypbind

. /etc/rc.subr
. /etc/rc.conf

name="amd"
command="/usr/sbin/${name}"
command_args='-p -a '$amd_dir' `sed s/#.*$// <'$amd_master'` >/var/run/amd.pid'
required_files="$amd_master"
required_dirs="$amd_dir"
required_vars="portmap"

run_rc_command "$1"
