# etc/rc.d/autounmountd
#!/bin/sh
#
# $NetBSD: autounmountd,v 1.1.4.1 2020/04/08 14:03:57 martin Exp $
#

# PROVIDE: autounmountd
# REQUIRE: DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="autounmountd"
rcvar="autofs"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
