# etc/rc.d/bootparams
#!/bin/sh
#
# $NetBSD: bootparams,v 1.1.1.1 2000/03/10 11:53:25 lukem Exp $
#

# PROVIDE: bootparams
# REQUIRE: portmap daemon

. /etc/rc.subr
. /etc/rc.conf

name="bootparamd"
command="/usr/sbin/rpc.${name}"
required_files="/etc/bootparams"

run_rc_command "$1"
