# etc/rc.d/bootparams
#!/bin/sh
#
# $NetBSD: bootparams,v 1.3 2000/05/13 08:45:06 lukem Exp $
#

# PROVIDE: bootparams
# REQUIRE: portmap DAEMON

. /etc/rc.subr

name="bootparamd"
command="/usr/sbin/rpc.${name}"
required_files="/etc/bootparams"

load_rc_config $name
run_rc_command "$1"
