# etc/rc.d/ccd
#!/bin/sh
#
# $NetBSD: ccd,v 1.4 2003/03/19 08:54:00 wiz Exp $
#

# PROVIDE: disks

. /etc/rc.subr

name="ccd"
rcvar=$name
start_cmd="ccd_start"
stop_cmd=":"

ccd_start()
{
	if [ -f /etc/ccd.conf ]; then
		echo "Configuring CCD devices."
		ccdconfig -C
	fi
}

load_rc_config $name
run_rc_command "$1"
