# etc/rc.d/cgd
#!/bin/sh
#
# $NetBSD: cgd,v 1.2 2003/03/19 06:06:48 wiz Exp $
#

# PROVIDE: disks

. /etc/rc.subr

name="cgd"
rcvar=$name
start_cmd="cgd_start"
stop_cmd=":"

cgd_start()
{
	if [ -f /etc/cgd/cgd.conf ]; then
		echo "Configuring CGD devices."
		cgdconfig -C
	fi
}

load_rc_config $name
run_rc_command "$1"
