# etc/rc.d/cron
#!/bin/sh
#
# $NetBSD: cron,v 1.5 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: cron
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="cron"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
