# etc/rc.d/cron
#!/bin/sh
#
# $NetBSD: cron,v 1.1.1.1 2000/03/10 11:53:25 lukem Exp $
#

# PROVIDE: cron
# REQUIRE: login

. /etc/rc.subr
. /etc/rc.conf

name="cron"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

run_rc_command "$1"
