# etc/rc.d/cron
#!/bin/sh
#
# $NetBSD: cron,v 1.2 2000/03/13 04:03:55 lukem Exp $
#

# PROVIDE: cron
# REQUIRE: LOGIN

. /etc/rc.subr
. /etc/rc.conf

name="cron"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

run_rc_command "$1"
