# etc/rc.d/dhclient
#!/bin/sh
#
# $NetBSD: dhclient,v 1.5 2000/05/26 20:44:17 hubertf Exp $
#

# PROVIDE: dhclient
# REQUIRE: systemfs syslogd

. /etc/rc.subr

name="dhclient"
command="/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
