# etc/rc.d/dhclient
#!/bin/sh
#
# $NetBSD: dhclient,v 1.3 2000/04/06 17:19:13 thorpej Exp $
#

# PROVIDE: dhclient
# REQUIRE: systemfs syslogd

. /etc/rc.subr
. /etc/rc.conf

name="dhclient"
command="/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1"
