# etc/rc.d/dhclient
#!/bin/sh
#
# $NetBSD: dhclient,v 1.2 2000/03/15 03:15:47 fvdl Exp $
#

# PROVIDE: dhclient
# REQUIRE: systemfs syslogd

. /etc/rc.subr
. /etc/rc.conf

name="dhclient"
command="/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf /var/db/${name}.leases"

run_rc_command "$1"
