# etc/rc.d/dhcpd
#!/bin/sh
#
# $NetBSD: dhcpd,v 1.1.1.1 2000/03/10 11:53:25 lukem Exp $
#

# PROVIDE: dhcpd
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="dhcpd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf /var/db/${name}.leases"

run_rc_command "$1"
