# etc/rc.d/dhcrelay
#!/bin/sh
#
# $NetBSD: dhcrelay,v 1.1.1.1 2000/03/10 11:53:25 lukem Exp $
#

# PROVIDE: dhcrelay
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="dhcrelay"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

run_rc_command "$1"
