# etc/rc.d/dmesg
#!/bin/sh
#
# $NetBSD: dmesg,v 1.6 2001/04/25 23:23:37 lukem Exp $
#

# PROVIDE: dmesg
# REQUIRE: mountcritremote

. /etc/rc.subr

name="dmesg"
rcvar=$name
dmesg_file="/var/run/dmesg.boot"
start_cmd="do_dmesg"
stop_cmd=":"

do_dmesg()
{
	rm -f ${dmesg_file}
	( umask 022 ; /sbin/dmesg $dmesg_flags > ${dmesg_file} )
}

load_rc_config $name
run_rc_command "$1"
