# etc/rc.d/gated
#!/bin/sh
#
# $NetBSD: gated,v 1.1.1.1 2000/03/10 11:53:25 lukem Exp $
#

# PROVIDE: gated
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="gated"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1"
