# etc/rc.d/gated
#!/bin/sh
#
# $NetBSD: gated,v 1.2 2000/03/13 04:03:57 lukem Exp $
#

# PROVIDE: gated
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="gated"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1"
