# etc/rc.d/ifwatchd
#!/bin/sh
#
# $NetBSD: ifwatchd,v 1.1 2001/12/10 16:53:07 martin Exp $
#

# PROVIDE: ifwatchd
# REQUIRE: mountcritremote network

. /etc/rc.subr

name="ifwatchd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"

