# etc/rc.d/ifwatchd
#!/bin/sh
#
# $NetBSD: ifwatchd,v 1.2 2003/03/14 09:18:50 martin Exp $
#

# PROVIDE: ifwatchd
# REQUIRE: mountcritremote
# BEFORE: SERVERS

. /etc/rc.subr

name="ifwatchd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"

