# etc/rc.d/inetd
#!/bin/sh
#
# $NetBSD: inetd,v 1.4 2000/05/13 08:45:06 lukem Exp $
#

# PROVIDE: inetd
# REQUIRE: DAEMON LOGIN

. /etc/rc.subr

name="inetd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
