# etc/rc.d/inetd
#!/bin/sh
#
# $NetBSD: inetd,v 1.4.4.1 2000/08/09 19:09:40 lukem Exp $
#

# PROVIDE: inetd
# REQUIRE: DAEMON LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="inetd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
