# etc/rc.d/inetd
#!/bin/sh
#
# $NetBSD: inetd,v 1.2 2000/03/13 04:03:57 lukem Exp $
#

# PROVIDE: inetd
# REQUIRE: DAEMON LOGIN

. /etc/rc.subr
. /etc/rc.conf

name="inetd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1" "reload"
