# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.9 2002/04/18 05:02:01 lukem Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd
# BEFORE:  SERVERS

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"
command_args="-D"

load_rc_config $name
run_rc_command "$1"
