# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.4 2001/07/05 05:26:58 tron Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
