# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.10.2.1 2004/04/06 03:56:14 jmc Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd
# BEFORE:  SERVERS

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
