# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.11 2004/04/04 13:59:36 mrg Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd
# BEFORE:  SERVERS

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
