# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.5 2002/02/28 01:06:13 lukem Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"
command_args="-D"

load_rc_config $name
run_rc_command "$1"
