# etc/rc.d/kerberos
#!/bin/sh
#
# $NetBSD: kerberos,v 1.4 2000/05/13 08:45:07 lukem Exp $
#

# PROVIDE: kerberos
# REQUIRE: DAEMON LOGIN

. /etc/rc.subr

arg=$1

name="kerberos"
command="/usr/sbin/${name}"
command_args=">> /var/log/${name}.log &"

load_rc_config $name
run_rc_command "$arg"


name="kadmind"
rcvar="kerberos"
command="/usr/sbin/${name}"
command_args="-n >> /var/log/${name}.log &"

# Note: config is shared with kerberos
run_rc_command "$arg"
