# etc/rc.d/kerberos
#!/bin/sh
#
# $NetBSD: kerberos,v 1.2 2000/03/13 04:03:57 lukem Exp $
#

# PROVIDE: kerberos
# REQUIRE: DAEMON LOGIN

. /etc/rc.subr
. /etc/rc.conf

name="kerberos"
command="/usr/sbin/${name}"
command_args=">> /var/log/${name}.log &"

run_rc_command "$1"


name="kadmind"
rcvar="kerberos"
command="/usr/sbin/${name}"
command_args="-n >> /var/log/${name}.log &"

run_rc_command "$1"
