# etc/rc.d/ldconfig
#!/bin/sh
#
# $NetBSD: ldconfig,v 1.3 2000/09/29 17:49:25 fvdl Exp $
#

# PROVIDE: ldconfig
# REQUIRE: mountall

. /etc/rc.subr

name="ldconfig"
start_cmd="ldconfig_start"
stop_cmd=":"

ldconfig_start()
{
	if [ -f /sbin/ldconfig ]; then
		echo "Creating a.out runtime link editor directory cache."
		ldconfig
	fi
}

load_rc_config $name
run_rc_command "$1"
