# etc/rc.d/ldconfig
#!/bin/sh
#
# $NetBSD: ldconfig,v 1.6.2.1 2004/06/07 10:28:22 tron Exp $
#

# PROVIDE: ldconfig
# REQUIRE: mountall
# BEFORE:  DAEMON

. /etc/rc.subr

name="ldconfig"
rcvar="ldconfig"
ldconfig_command="/sbin/ldconfig"
start_cmd="ldconfig_start"
stop_cmd=":"

ldconfig_start()
{
	if [ -f ${ldconfig_command} ]; then
		echo "Creating a.out runtime link editor directory cache."
		${ldconfig_command} ${rc_flags}
	fi
}

load_rc_config $name
run_rc_command "$1"
