# etc/rc.d/lkm3
#!/bin/sh
#
# $NetBSD: lkm3,v 1.2 2000/03/11 20:10:21 veego Exp $
#

# PROVIDE: aftermountlkm
# REQUIRE: mountall

. /etc/rc.subr
. /etc/rc.conf

name="lkm3"
start_precmd="checkyesno lkm"
start_cmd="do_lkm3 start"
stop_precmd="$start_precmd"
stop_cmd="do_lkm3 stop"

do_lkm3()
{
	#	(un)load kernel modules specified in /etc/lkm.conf
	#
	if checkyesno lkm && [ -f /etc/rc.lkm ]; then
		lkmstage=AFTERMOUNT
		set $1 ; . /etc/rc.lkm
	fi
}

run_rc_command "$1"
