# etc/rc.d/local
#!/bin/sh
#
# $NetBSD: local,v 1.4 2000/05/13 08:45:07 lukem Exp $
#

# REQUIRE: DAEMON
# PROVIDE: local
# BEFORE: LOGIN

. /etc/rc.subr

name="local"
start_cmd="do_local"
stop_cmd=":"

do_local()
{
	if [ -f /etc/rc.local ]; then
		. /etc/rc.local
	fi
}

load_rc_config $name
run_rc_command "$1"
