# etc/rc.d/lpd
#!/bin/sh
#
# $NetBSD: lpd,v 1.2 2000/03/13 04:04:00 lukem Exp $
#

# PROVIDE: lpd
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="lpd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/printcap"

run_rc_command "$1"
