# etc/rc.d/mixerctl
#!/bin/sh
#
# $NetBSD: mixerctl,v 1.4 2002/07/16 03:24:24 rafal Exp $
#

# PROVIDE: mixerctl
# REQUIRE: mountcritremote

. /etc/rc.subr

name="mixerctl"
start_cmd="mixerctl_start"
stop_cmd=":"

mixerctl_start()
{
	if [ -r /etc/mixerctl.conf ]; then
		echo "Setting mixerctl variables..."
		while read setting; do
			[ -z "$setting" ] || mixerctl -n -w $setting
		done < /etc/mixerctl.conf
	fi
}

load_rc_config $name
run_rc_command "$1"
