# etc/rc.d/motd
#!/bin/sh
#
# $NetBSD: motd,v 1.2.4.1 2000/08/01 13:42:25 chuck Exp $
#

# PROVIDE: motd
# REQUIRE: systemfs

. /etc/rc.subr

name="motd"
start_precmd="checkyesno update_motd"
start_cmd="motd_start"
stop_cmd=":"

motd_start()
{
	#	Update kernel info in /etc/motd
	#	Must be done *before* interactive logins are possible
	#	to prevent possible race conditions.
	#
	echo "Updating motd."
	if [ ! -f /etc/motd ]; then
		install -c -o root -g wheel -m 664 /dev/null /etc/motd
	fi
	T=/etc/_motd
	rm -f $T
	sysctl -n kern.version | sed 1q > $T
	sed '1{/^NetBSD.*/{d;};};' < /etc/motd >> $T
	cmp -s $T /etc/motd || cp $T /etc/motd
	rm -f $T
}

load_rc_config $name
run_rc_command "$1"
