# etc/rc.d/mountcritremote
#!/bin/sh
#
# $NetBSD: mountcritremote,v 1.1.1.1 2000/03/10 11:53:24 lukem Exp $
#

# PROVIDE: mountcritremote
# REQUIRE: root mountcritlocal network

. /etc/rc.subr
. /etc/rc.conf

name="mountcritremote"
start_cmd="mount_critical_filesystems remote"
stop_cmd=":"

run_rc_command "$1"
