# etc/rc.d/mountd
#!/bin/sh
#
# $NetBSD: mountd,v 1.4.2.1 2000/06/23 15:46:03 minoura Exp $
#

# PROVIDE: mountd
# REQUIRE: mountall beforemountlkm network rpcbind quota

. /etc/rc.subr

name="mountd"
rcvar="nfs_server"
command="/usr/sbin/${name}"
required_files="/etc/exports"
start_precmd="mountd_precmd"
extra_commands="reload"

mountd_precmd()
{
	rm -f /var/db/mountdtab
	echo -n > /var/db/mountdtab
	return 0
}

load_rc_config $name
run_rc_command "$1"
