# etc/rc.d/mrouted
#!/bin/sh
#
# $NetBSD: mrouted,v 1.2 2000/03/13 04:04:01 lukem Exp $
#

# PROVIDE: mrouted
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="mrouted"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1" "reload"
