# etc/rc.d/mrouted
#!/bin/sh
#
# $NetBSD: mrouted,v 1.1.1.1 2000/03/10 11:53:24 lukem Exp $
#

# PROVIDE: mrouted
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="mrouted"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

run_rc_command "$1" "reload"
