# etc/rc.d/named
#!/bin/sh
#
# $NetBSD: named,v 1.2 2000/03/13 04:04:01 lukem Exp $
#

# PROVIDE: named
# REQUIRE: SERVERS

. /etc/rc.subr
. /etc/rc.conf

name="named"
command="/usr/sbin/${name}"
command_flags="$1"
pidfile="/var/run/${name}.pid"

run_rc_command "$1" "reload dumpdb notrace qrylog querylog stats trace"
