# etc/rc.d/named
#!/bin/sh
#
# $NetBSD: named,v 1.4 2000/05/13 03:46:42 lukem Exp $
#

# PROVIDE: named
# REQUIRE: SERVERS

. /etc/rc.subr
. /etc/rc.conf

name="named"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

extra_commands="reload"

run_rc_command "$1"
