# etc/rc.d/nfslocking
#!/bin/sh
#
# $NetBSD: nfslocking,v 1.2 2000/04/15 21:14:50 tsarna Exp $
#

# PROVIDE: nfslocking
# REQUIRE: nfsd

. /etc/rc.subr
. /etc/rc.conf

arg=$1

start_precmd='checkyesno nfs_server || checkyesno nfs_client'
stop_precmd=$start_precmd
status_precmd=$start_precmd

name="statd"
command="/usr/sbin/rpc.${name}"
run_rc_command "$arg"

name="lockd"
command="/usr/sbin/rpc.${name}"
run_rc_command "$arg"
