# etc/rc.d/nfslocking
#!/bin/sh
#
# $NetBSD: nfslocking,v 1.4 2000/09/19 13:04:39 lukem Exp $
#

# PROVIDE: nfslocking
# REQUIRE: nfsd

. /etc/rc.subr

arg=$1

start_precmd='checkyesno nfs_server || checkyesno nfs_client'
stop_precmd=$start_precmd
status_precmd=$start_precmd

name="statd"
rcvar=$name
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"

name="lockd"
rcvar=$name
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"
