# etc/rc.d/ntpd
#!/bin/sh
#
# $NetBSD: ntpd,v 1.1 2000/05/02 12:19:20 simonb Exp $
#

# PROVIDE: ntpd
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="ntpd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ntp.conf"

run_rc_command "$1"
