# etc/rc.d/ntpdate
#!/bin/sh
#
# $NetBSD: ntpdate,v 1.2 2000/03/11 20:10:21 veego Exp $
#

# PROVIDE: ntpdate
# REQUIRE: syslogd network

. /etc/rc.subr
. /etc/rc.conf

name="ntpdate"
start_precmd="checkyesno ntpdate"
start_cmd="ntpdate_start"
stop_cmd=":"

ntpdate_start()
{
	if [ -z "$ntpdate_hosts" ]; then
		ntpdate_hosts=`awk '
		    /^server[ \t]*127.127/	{next}
		    /^(server|peer)/		{print $2}
		' </etc/ntp.conf`
	fi
	if [ -n "$ntpdate_hosts"  ]; then
		echo "Setting date via ntp."
		ntpdate -b -s $ntpdate_hosts
	fi
}

run_rc_command "$1"
