# etc/rc.d/ntpdate
#!/bin/sh
#
# $NetBSD: ntpdate,v 1.9 2002/07/19 03:22:09 bad Exp $
#

# PROVIDE: ntpdate
# REQUIRE: NETWORKING syslogd

. /etc/rc.subr

name="ntpdate"
rcvar=$name
start_cmd="ntpdate_start"
stop_cmd=":"

ntpdate_start()
{
	if [ -z "$ntpdate_hosts" ]; then
		ntpdate_hosts=`awk '
		    /^(server|peer)[ \t]*127.127/	{next}
		    /^(server|peer)/		{print $2}
		' </etc/ntp.conf`
	fi
	if [ -n "$ntpdate_hosts"  ]; then
		echo "Setting date via ntp."
		ntpdate $rc_flags $ntpdate_hosts
	fi
}

load_rc_config $name
run_rc_command "$1"
