# etc/rc.d/portmap
#!/bin/sh
#
# $NetBSD: portmap,v 1.2 2000/05/13 08:45:08 lukem Exp $
#

# PROVIDE: portmap
# REQUIRE: network ntpdate syslogd named ppp

. /etc/rc.subr

name="portmap"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
