# etc/rc.d/postfix
#!/bin/sh
#
# $NetBSD: postfix,v 1.3 2000/04/30 12:21:00 lukem Exp $
#

# PROVIDE: mail
# REQUIRE: LOGIN

#	we could do this, but make mail start late, so that things like
#	.forward's are not processed until the system is fully operational
## REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="postfix"
required_files="/etc/${name}/main.cf"
extra_commands="reload"

start_precmd="checkyesno postfix"
start_cmd="postfix start"

stop_precmd=$start_precmd
stop_cmd="postfix stop"

reload_precmd=$start_precmd
reload_cmd="postfix reload"

run_rc_command "$1"
