# etc/rc.d/ppp
#!/bin/sh
#
# $NetBSD: ppp,v 1.2.4.1 2000/08/09 19:09:45 lukem Exp $
#

# PROVIDE: ppp
# REQUIRE: mountcritremote syslogd
#
#	Note that this means that syslogd will not be listening on
#	any PPP addresses.  This is considered a feature.
#

. /etc/rc.subr

name="ppp"
start_cmd="ppp_start"
stop_cmd=":"

ppp_start()
{
	#	/etc/ppp/peers and $ppp_peers contain boot configuration
	#	information for pppd.  each value in $ppp_peers that has a
	#	file in /etc/ppp/peers of the same name, will be run as
	#	`pppd call <peer>'.
	#
	if [ -n "$ppp_peers" ]; then
		set -- $ppp_peers
		echo -n "Starting pppd:"
		while [ $# -ge 1 ]; do
			peer=$1
			shift
			if [ -f /etc/ppp/peers/$peer ]; then
				pppd call $peer
				echo -n " $peer"
			fi
		done
		echo "."
	fi
}

load_rc_config $name
run_rc_command "$1"
