# etc/rc.d/raidframe
#!/bin/sh
#
# $NetBSD: raidframe,v 1.2 2000/04/26 05:13:51 thorpej Exp $
#

# PROVIDE: disks

. /etc/rc.subr
. /etc/rc.conf

name="raidframe"
start_cmd="raidframe_start"
stop_cmd=":"

raidframe_start()
{
	#	Configure raid devices.
	#
	for dev in 0 1 2 3; do
		if [ -f /etc/raid$dev.conf ]; then
			raidctl -c /etc/raid$dev.conf raid$dev
		fi
	done
}

run_rc_command "$1"
