# etc/rc.d/rarpd
#!/bin/sh
#
# $NetBSD: rarpd,v 1.1.1.1 2000/03/10 11:53:24 lukem Exp $
#

# PROVIDE: rarpd
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="rarpd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ethers"

run_rc_command "$1"
