# etc/rc.d/rbootd
#!/bin/sh
#
# $NetBSD: rbootd,v 1.3 2000/04/30 13:23:35 lukem Exp $
#

# PROVIDE: rbootd
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="rbootd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

run_rc_command "$1"
