# etc/rc.d/route6d
#!/bin/sh
#
# $NetBSD: route6d,v 1.2 2000/03/13 04:04:03 lukem Exp $
#

# PROVIDE: route6d
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="route6d"
command="/usr/sbin/${name}"

run_rc_command "$1"
