# etc/rc.d/route6d
#!/bin/sh
#
# $NetBSD: route6d,v 1.1.1.1 2000/03/10 11:53:24 lukem Exp $
#

# PROVIDE: route6d
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="route6d"
command="/usr/sbin/${name}"

run_rc_command "$1"
