# etc/rc.d/routed
#!/bin/sh
#
# $NetBSD: routed,v 1.4.4.1 2000/10/02 01:09:26 lukem Exp $
#

# PROVIDE: routed
# REQUIRE: DAEMON gated

. /etc/rc.subr

name="routed"
rcvar=$name
command="/sbin/${name}"
start_precmd="routed_precmd"

routed_precmd()
{
	if checkyesno gated && checkyesno routed; then
		warn "gated and routed both requested to be run: only running gated."
		return 1
	fi
}

load_rc_config $name
run_rc_command "$1"
