# etc/rc.d/rpcbind
#!/bin/sh
#
# $NetBSD: rpcbind,v 1.2 2000/06/03 00:17:17 fvdl Exp $
#

# PROVIDE: rpcbind
# REQUIRE: network ntpdate syslogd named ppp

. /etc/rc.subr
. /etc/rc.conf

name="rpcbind"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

run_rc_command "$1"
