# etc/rc.d/rtsold
#!/bin/sh
#
# $NetBSD: rtsold,v 1.2 2000/03/13 04:04:04 lukem Exp $
#

# PROVIDE: rtsold
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="rtsold"
command="/usr/sbin/${name}"
start_precmd="rtsold_precmd"

rtsold_precmd()
{
	if [ "$ip6mode" != "autohost" ]; then
		warn "\$ip6mode must be set to 'autohost' to use ${name}."
		return 1
	fi
}

run_rc_command "$1"
