# etc/rc.d/sendmail
#!/bin/sh
#
# $NetBSD: sendmail,v 1.2 2000/03/13 04:04:05 lukem Exp $
#

# PROVIDE: mail
# REQUIRE: LOGIN

#	we could do this, but make mail start late, so that things like
#	.forward's are not processed until the system is fully operational
## REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="sendmail"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.cf"

run_rc_command "$1"
