# etc/rc.d/sendmail
#!/bin/sh
#
# $NetBSD: sendmail,v 1.6 2000/05/15 06:16:04 tron Exp $
#

# PROVIDE: mail
# REQUIRE: LOGIN

#	we could do this, but make mail start late, so that things like
#	.forward's are not processed until the system is fully operational
## REQUIRE: DAEMON

. /etc/rc.subr

name="sendmail"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/mail/${name}.cf"

load_rc_config $name

# die if there's pre-8.10 configuration file
if checkyesno sendmail && [ -f "/etc/${name}.cf" ]; then
	err 1 "${name} was not started, move /etc/${name}.cf to /etc/mail."
fi

run_rc_command "$1"
