# etc/rc.d/sendmail
#!/bin/sh
#
# $NetBSD: sendmail,v 1.3 2000/05/03 10:55:20 itojun Exp $
#

# PROVIDE: mail
# REQUIRE: LOGIN

#	we could do this, but make mail start late, so that things like
#	.forward's are not processed until the system is fully operational
## REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="sendmail"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
if [ -f "/etc/${name}.cf" ]; then
	required_files="/etc/${name}.cf"
	command_args="-C${required_files}"
	warn "using ${required_files}."
	warn "default location was changed, move ${required_files} to /etc/mail."
else
	required_files="/etc/mail/${name}.cf"
fi

run_rc_command "$1"
