# etc/rc.d/swap1
#!/bin/sh
#
# $NetBSD: swap1,v 1.3 2000/04/26 05:13:51 thorpej Exp $
#

# PROVIDE: localswap
# REQUIRE: disks

. /etc/rc.subr
. /etc/rc.conf

name="swap1"

#		Add all block-type swap devices; these might be necessary
#		during disk checks.
#
start_cmd='swapctl -A -t blk'

#		Remove all block-type swap devices
#
stop_cmd="swapctl -U -t blk"

run_rc_command "$1"
