# etc/rc.d/swap1
#!/bin/sh
#
# $NetBSD: swap1,v 1.4.4.2 2001/09/30 15:50:15 he Exp $
#

# PROVIDE: localswap
# REQUIRE: disks

. /etc/rc.subr

name="swap1"
start_cmd="swap1_start"

#		Add all block-type swap devices; these might be necessary
#		during disk checks.
#
swap1_start()
{
	if ! checkyesno no_swap; then
		swapctl -A -t blk
	fi
}

#		Remove all block-type swap devices
#
stop_cmd="swapctl -U -t blk"

load_rc_config swap
run_rc_command "$1"
