# etc/rc.d/swap1
#!/bin/sh
#
# $NetBSD: swap1,v 1.2 2000/03/11 20:12:22 veego Exp $
#

# PROVIDE: localswap
# REQUIRE: ccd raidframe

. /etc/rc.subr
. /etc/rc.conf

name="swap1"

#		Add all block-type swap devices; these might be necessary
#		during disk checks.
#
start_cmd='swapctl -A -t blk'

#		Remove all block-type swap devices
#
stop_cmd="swapctl -U -t blk"

run_rc_command "$1"
