# etc/rc.d/sysctl
#!/bin/sh
#
# $NetBSD: sysctl,v 1.11 2002/03/27 08:03:08 lukem Exp $
#

# PROVIDE: sysctl
# REQUIRE: root ipfilter ipsec
# BEFORE:  DAEMON

. /etc/rc.subr

name="sysctl"
start_cmd="sysctl_start"
stop_cmd=":"

sysctl_start()
{
	if [ -n "$defcorename" ]; then
		warn 'rc.conf(5) $defcorename setting is obsolete.'
		warn 'Please set kern.defcorename in /etc/sysctl.conf'
	fi

	if [ -n "$nfsiod_flags" ]; then
		warn 'rc.conf(5) $nfsiod_flags setting is obsolete.'
		warn 'Please set vfs.nfs.iothreads in /etc/sysctl.conf'
	fi

	if [ -r /etc/sysctl.conf ]; then
		echo "Setting sysctl variables:"
		sysctl -f /etc/sysctl.conf
	fi
}

load_rc_config $name
run_rc_command "$1"
