# etc/rc.d/sysctl
#!/bin/sh
#
# $NetBSD: sysctl,v 1.2 2000/04/08 23:02:54 tsarna Exp $
#

# PROVIDE: sysctl
# REQUIRE: root ipnat

. /etc/rc.subr
. /etc/rc.conf

name="sysctl"
start_cmd="sysctl_start"
stop_cmd=":"

sysctl_start()
{
	if [ -r /etc/sysctl.conf ]; then
		echo "Setting sysctl variables:"
		sysctl -f /etc/sysctl.conf
	fi
}

run_rc_command "$1"
