# etc/rc.d/timed
#!/bin/sh
#
# $NetBSD: timed,v 1.1.1.1 2000/03/10 11:53:24 lukem Exp $
#

# PROVIDE: timed
# REQUIRE: daemon

. /etc/rc.subr
. /etc/rc.conf

name="timed"
command="/usr/sbin/${name}"

run_rc_command "$1"
