# etc/rc.d/wscons
#!/bin/sh
#
# $NetBSD: wscons,v 1.3 2000/05/13 08:45:10 lukem Exp $
#

# PROVIDE: wscons
# REQUIRE: systemfs

. /etc/rc.subr

name="wscons"
start_precmd="checkyesno wscons"
start_cmd="do_wscons"
stop_cmd=":"

do_wscons()
{
	if [ -f /etc/rc.wscons ]; then
		. /etc/rc.wscons
	fi
}

load_rc_config $name
run_rc_command "$1"
