# etc/rc.d/wsmoused
#!/bin/sh
#
# $NetBSD: wsmoused,v 1.2 2003/03/04 22:32:11 jmmv Exp $
#

# PROVIDE: wsmoused
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="wsmoused"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
