# etc/rc.d/xdm
#!/bin/sh
#
# $NetBSD: xdm,v 1.4.4.1 2000/08/09 19:09:50 lukem Exp $
#

# PROVIDE: xdm
# REQUIRE: DAEMON LOGIN wscons
# KEYWORD: shutdown

. /etc/rc.subr

name="xdm"
command="/usr/X11R6/bin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/usr/X11R6/lib/X11/xdm/xdm-config"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
