# etc/rc.d/xfs
#!/bin/sh
#
# $NetBSD: xfs,v 1.2 2000/04/30 13:23:04 lukem Exp $
#

# PROVIDE: xfs
# REQUIRE: mountall

. /etc/rc.subr
. /etc/rc.conf

name="xfs"
command="/usr/X11R6/bin/${name}"
command_args="& sleep 2"
required_files="/usr/X11R6/lib/X11/fs/config"
extra_commands="reload"

run_rc_command "$1"
