# etc/rc.d/xfs
#!/bin/sh
#
# $NetBSD: xfs,v 1.3.4.1 2000/10/02 01:08:58 lukem Exp $
#

# PROVIDE: xfs
# REQUIRE: mountall

. /etc/rc.subr

name="xfs"
rcvar=$name
command="/usr/X11R6/bin/${name}"
command_args="& sleep 2"
required_files="/usr/X11R6/lib/X11/fs/config"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
