# etc/rc.d/xntpd
#!/bin/sh
#
# $NetBSD: xntpd,v 1.2 2000/03/13 04:04:11 lukem Exp $
#

# PROVIDE: xntpd
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/rc.conf

name="xntpd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ntp.conf"

run_rc_command "$1"
