# etc/rc.d/ypbind
#!/bin/sh
#
# $NetBSD: ypbind,v 1.2 2000/03/16 22:22:15 lukem Exp $
#

# PROVIDE: ypbind
# REQUIRE: ypserv

. /etc/rc.subr
. /etc/rc.conf

name="ypbind"
command="/usr/sbin/${name}"
start_precmd="ypbind_precmd"

ypbind_precmd()
{
	_domain=`domainname`
	if [ -z "$_domain" ]; then
		warn "domainname(1) is not set; ${name} not started"
		return 1
	fi
}

run_rc_command "$1"
