# etc/rc.d/yppasswdd
#!/bin/sh
#
# $NetBSD: yppasswdd,v 1.2 2000/03/16 22:22:15 lukem Exp $
#

# PROVIDE: yppasswdd
# REQUIRE: ypserv ypbind

. /etc/rc.subr
. /etc/rc.conf

name="yppasswdd"
command="/usr/sbin/rpc.${name}"
required_vars="portmap ypserv"
start_precmd="yppasswdd_precmd"

yppasswdd_precmd()
{
	_domain=`domainname`
	if [ -z "$_domain" ]; then
		warn "domainname(1) is not set; ${name} not started"
		return 1
	fi
}

run_rc_command "$1"
